# deploy/systemd/talisman-gateway.service
[Unit]
Description=TalisMan host-side credential and budget gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/TalisMan
ExecStart=%h/TalisMan/.venv-gateway/bin/uvicorn TalisMan_gateway:app --app-dir %h/TalisMan/src --host 127.0.0.1 --port 8787
Restart=on-failure
RestartSec=5
Environment=TALISMAN_HOME=%h/TalisMan

# Secrets are file-based and not committed to Git. For stronger hardening, migrate these to systemd LoadCredential.
ReadWritePaths=%h/TalisMan/state %h/TalisMan/logs
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
